# crontab
*/10 * * * * root /godaddy-ddns.sh >> /var/log/cron.log 2>&1

