# crontab
*/15 * * * * /bin/bash /godaddy-ddns.sh

