# crontab
*/10 * * * * root /godaddy-ddns.sh > /var/log/godaddy-ddns.log 2>&1

